# layout/Library/LaunchDaemons/com.jb-detect-bypass.applaunchmond.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>Label</key>
	<string>com.jb-detect-bypass.applaunchmond</string>
	<key>ProgramArguments</key>
	<array>
		<string>/var/jb/usr/local/libexec/applaunchmond</string>
		<string>--launchd</string>
	</array>
	<key>UserName</key>
	<string>root</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
</dict>
</plist>
